# deploy/systemd/system/bluealsa.service.d/override.conf
# TySpeaker streams Bluetooth A2DP to the speaker via bluez-alsa. The Pi is the
# A2DP *source* (it sends audio to the speaker), so run source-only.
#
# --keep-alive=600 holds the A2DP transport open for 600s after playback so the
# speaker does NOT idle-sleep and drop the connection between sounds. This is the
# fix for the "connects then drops in ~20s" behaviour every speaker showed.
[Service]
ExecStart=
ExecStart=/usr/bin/bluealsa --keep-alive=600 -p a2dp-source
