# deploy/tyspeaker.service
[Unit]
Description=TySpeaker push-button audio player + web UI
After=network-online.target sound.target bluetooth.target bluealsa.service
Wants=network-online.target bluealsa.service

[Service]
Type=simple
User=jaysettle
# Adjust if you cloned/installed elsewhere:
WorkingDirectory=/home/jaysettle/tyspeaker-app
Environment=TYSPEAKER_DATA=/home/jaysettle/tyspeaker
Environment=TYSPEAKER_PORT=8000
Environment=TYSPEAKER_BUTTON_PIN=17
Environment="TYSPEAKER_PLAYER=mpv --no-video --really-quiet --audio-device=alsa/bluealsa {path}"
ExecStart=/home/jaysettle/tyspeaker-app/.venv/bin/python -m tyspeaker
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
