# deploy/systemd/system/bluetooth.service.d/10-tyspeaker-clean-cache.conf
[Service]
Restart=on-failure
RestartSec=3
ExecStartPre=/usr/local/sbin/tyspeaker-clean-bluez-lastused
