# deploy/systemd/system/tyspeaker.service.d/offline.conf
# TySpeaker must boot and play with NO network (the scooter leaves home WiFi).
# The stock unit has After=/Wants=network-online.target, which — together with an
# enabled NetworkManager-wait-online.service — stalls boot for ~90s waiting for a
# network that will never appear when away from home. The app needs only local
# files + Bluetooth, so drop the network dependency entirely.
#
# Also run: sudo systemctl disable NetworkManager-wait-online.service
[Unit]
After=
After=sound.target bluetooth.service bluealsa.service
Wants=
